# telegram-voting-app/telegram-voting-app.service
[Unit]
Description=Telegram Voting App
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt
ExecStart=/usr/bin/npm run preview
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target